# wazuh-rules/local_rules_combined.xml
<!--
Combined Wazuh local rules for the Windows SIEM Detection Lab.
These rules are for lab demonstration and should be tuned before production use.
-->

<group name="windows,sysmon,attack,custom_lab,">

  <rule id="100100" level="10">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)\\(powershell|pwsh)\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(-enc|-encodedcommand|frombase64string)</field>
    <description>Suspicious PowerShell encoded command execution detected</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>windows,powershell,execution,attack,</group>
  </rule>

  <rule id="100101" level="9">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)\\schtasks\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(/create|/sc|/tn|/tr)</field>
    <description>Scheduled task creation or modification detected</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>windows,scheduled_task,persistence,attack,</group>
  </rule>

  <rule id="100102" level="10">
    <if_group>windows_security</if_group>
    <field name="win.system.eventID">4720</field>
    <description>Local user account creation detected</description>
    <mitre>
      <id>T1136.001</id>
    </mitre>
    <group>windows,account_creation,persistence,attack,</group>
  </rule>

  <rule id="100112" level="9">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)\\net1?\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(user\s+\S+\s+\S+\s+/add|localgroup\s+administrators)</field>
    <description>Local user account creation command detected</description>
    <mitre>
      <id>T1136.001</id>
    </mitre>
    <group>windows,account_creation,persistence,attack,</group>
  </rule>

  <rule id="100103" level="10">
    <if_group>windows_system</if_group>
    <field name="win.system.eventID">7045</field>
    <description>Windows service creation detected</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
    <group>windows,service_creation,persistence,attack,</group>
  </rule>

  <rule id="100113" level="9">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)\\sc\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(create|config|binpath=)</field>
    <description>Windows service creation or modification command detected</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
    <group>windows,service_creation,persistence,attack,</group>
  </rule>

</group>